# bin/netdisk-fast-download.service
[Unit]
Description=netdisk-fast-download
Documentation=https://qaiu.top
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User=USER
# 需要JDK17及以上版本 注意修改为自己的路径
ExecStart=/root/java/jdk-17.0.2/bin/java -server -Xmx128m -jar /root/java/netdisk-fast-download-0.1.6/netdisk-fast-download-0.1.6.jar
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=always
StandOutput=syslog

StandError=inherit

[Install]
WantedBy=multi-user.target
